# data/group-admin-daemon.service.in
[Unit]
Description=Group Service

After=nss-user-lookup.target
Wants=nss-user-lookup.target

[Service]
Type=dbus
BusName=org.group.admin
ExecStart=@libexecdir@/group-admin-daemon
Environment=GVFS_DISABLE_FUSE=1
Environment=GIO_USE_VFS=local
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1

[Install]
WantedBy=graphical.target
